# infra/mac-relay/com.formicanera.pf-nat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.formicanera.pf-nat</string>
  <key>RunAtLoad</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>pfctl -e -f /etc/pf-formicanera.conf 2&gt;/dev/null; true</string>
  </array>
  <key>StandardOutPath</key>
  <string>/var/log/formicanera-pf-nat.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/formicanera-pf-nat.log</string>
</dict>
</plist>
